# backup/hostapd/hostapd.service
[Service]
ExecStartPre=/sbin/iw dev wlan0 interface add ap0 type __ap
